# pub/hi3716mv300/rootbox/etc/udev/rules.d/50-firmware.rules
# Load firmware
SUBSYSTEM=="firmware", ACTION=="add", RUN+="firmware.sh"
